# docker/backup/scripts/crontab
# Запуск бэкапа каждые 6 часов (00:00, 06:00, 12:00, 18:00 MSK)
0 */6 * * * /scripts/backup.sh >> /var/log/backup.log 2>&1
